# scheduler-app/crontab
* * * * * root bash charge-all-Invoices.sh
#
